# ssvm-backend/supervisord.conf
; ═══════════════════════════════════════════════════════════════
; SSVM Backend — Supervisord Configuration
; Manages both PostgreSQL and Express.js in a single container
; ═══════════════════════════════════════════════════════════════

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; ── PostgreSQL Database Process ──
[program:postgresql]
command=/usr/libexec/postgresql16/postgres -D /var/lib/postgresql/data
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ── Express.js API Server Process ──
[program:express]
command=node /app/dist/index.js
directory=/app
autostart=true
autorestart=true
priority=20
startsecs=10
startretries=5
environment=NODE_ENV="production",PORT="5000",DATABASE_URL="postgresql://%(ENV_POSTGRES_USER)s:%(ENV_POSTGRES_PASSWORD)s@localhost:5432/%(ENV_POSTGRES_DB)s",JWT_SECRET="%(ENV_JWT_SECRET)s",JWT_EXPIRES_IN="%(ENV_JWT_EXPIRES_IN)s",FRONTEND_URL="%(ENV_FRONTEND_URL)s",UPLOAD_DIR="%(ENV_UPLOAD_DIR)s"
stdout_logfile=/var/log/supervisor/express.log
stderr_logfile=/var/log/supervisor/express-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
